# internal/udev/70-intuneme-video.rules
# intuneme: forward video capture devices to systemd-nspawn container.
# Matches V4L2 video devices (/dev/video*) and media controllers (/dev/media*).
# Managed by intuneme — do not edit manually.

# V4L2 video capture devices (/dev/video*)
ACTION=="add", SUBSYSTEM=="video4linux", KERNEL=="video*", RUN+="/usr/local/lib/intuneme/usb-hotplug add /dev/%k"
ACTION=="remove", SUBSYSTEM=="video4linux", KERNEL=="video*", RUN+="/usr/local/lib/intuneme/usb-hotplug remove /dev/%k"

# Media controller devices (/dev/media*)
ACTION=="add", SUBSYSTEM=="media", KERNEL=="media*", RUN+="/usr/local/lib/intuneme/usb-hotplug add /dev/%k"
ACTION=="remove", SUBSYSTEM=="media", KERNEL=="media*", RUN+="/usr/local/lib/intuneme/usb-hotplug remove /dev/%k"
